# deploy/puzzpool-test.service
# puzzpool TEST instance — systemd service unit
#
# Install:
#   sudo cp deploy/puzzpool-test.service /etc/systemd/system/puzzpool-test.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now puzzpool-test
#
# Manage:
#   sudo systemctl status  puzzpool-test
#   sudo systemctl restart puzzpool-test
#   journalctl -u puzzpool-test -f

[Unit]
Description=puzzpool TEST — feature branch test instance
After=network.target

[Service]
Type=simple

User=jan
WorkingDirectory=/home/jan/git/puzzpool.test

ExecStart=/home/jan/git/puzzpool.test/build/bin/puzzpool

EnvironmentFile=/home/jan/git/puzzpool.test/.env
Environment=PORT=9999
Environment=DB_PATH=/home/jan/git/puzzpool.test/pool.db
Environment=TARGET_MINUTES=5
Environment=TIMEOUT_MINUTES=15
# Uncomment to enable admin token auth on the test instance:
# Environment=ADMIN_TOKEN=your-test-token-here

Restart=on-failure
RestartSec=5s
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
